# scripts/deploy/childnotes-api.service
[Unit]
Description=ChildNotes API (.NET 10)
Documentation=https://github.com/KleinPan/ChildNotes
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=childnotes
Group=childnotes
WorkingDirectory=/opt/childnotes/current
ExecStart=/opt/childnotes/current/ChildNotes.Api
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=childnotes-api
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:8080
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
Environment=DOTNET_NOLOGO=1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/childnotes /var/log/childnotes

[Install]
WantedBy=multi-user.target
